# infrastructure/ansible/files/datomic-upstart.conf
start on runlevel [2345]
 
pre-start script
bash << "EOF"
mkdir -p /var/log/datomic
chown -R datomic /var/log/datomic
EOF
end script
 
start on (started network-interface
or started network-manager
or started networking)
 
stop on (stopping network-interface
or stopping network-manager
or stopping networking)
 
respawn
 
script
exec su - datomic -c 'cd /var/lib/datomic/runtime; /var/lib/datomic/runtime/bin/transactor -Xms400m -Xmx400m /var/lib/datomic/transactor.properties 2>&1 >> /var/log/datomic/datomic.log'
end script
 
stop on runlevel [016]