# infrastructure/ansible/files/gp2-upstart.conf
start on started

script
PORT=80 /usr/bin/java -jar /var/projects/gp2/code/gp2.jar server >> /var/log/gp2/gp2.log 2>&1
end script
